# infrastructure/data/l2_ingestion/deploy/compress.service
# systemd unit (oneshot) — one compression pass (completed JSONL.gz -> Parquet).
# Triggered by compress.timer hourly; enable compress.timer, not this directly.

[Unit]
Description=Epsilon L2 compression pass (JSONL.gz -> Parquet)
# No network needed — purely local file conversion.

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/epsilon/l2_ingestion
# No --input: scan data/raw for completed (>1h old), not-yet-processed shards.
ExecStart=/opt/epsilon/l2_ingestion/venv/bin/python compression/pipeline.py

# Memory cap: each shard can peak ~1-2 GB (esports price_change has 3M+ rows).
# 3G leaves headroom for Python overhead while still protecting the 4 GB server.
MemoryMax=3G

StandardOutput=journal
StandardError=journal
SyslogIdentifier=l2-compress
