# infrastructure/data/l2_ingestion/deploy/sync.service
# systemd unit (oneshot) — one cloud-sync pass (parquet -> R2, raw backup, prune).
# Triggered by sync.timer every 6 hours; enable sync.timer, not this directly.

[Unit]
Description=Epsilon L2 cloud sync (rclone -> Cloudflare R2)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/epsilon/l2_ingestion
# The script resolves its own paths and assumes an rclone remote named "r2"
# is already configured (see DEPLOY.md). It exits nonzero on rclone failure,
# which systemd records as a failed unit.
ExecStart=/usr/bin/bash /opt/epsilon/l2_ingestion/sync/sync_cloud.sh

StandardOutput=journal
StandardError=journal
SyslogIdentifier=l2-sync
