# infrastructure/data/l2_ingestion/deploy/sync.timer
# systemd timer — runs sync.service every 6 hours.
# Install: copy to /etc/systemd/system/, then:
#   systemctl daemon-reload && systemctl enable --now sync.timer

[Unit]
Description=Run Epsilon L2 cloud sync every 6 hours

[Timer]
OnBootSec=20min
OnUnitActiveSec=6h
AccuracySec=2min
Unit=sync.service

[Install]
WantedBy=timers.target
